# deploy/spybot.service
[Unit]
Description=Spybot Next.js
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/Spybot
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/Spybot/.env
ExecStart=/usr/bin/pnpm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
